# crontab
PATH="/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/bin:/usr/local/sbin"
* * * * * npm start 2>&1 >> parser.log
